# deploy/frogtalk.service
[Unit]
Description=FrogChat standalone chat server
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/frogtalk
EnvironmentFile=/opt/frogtalk/.env
ExecStart=/opt/frogtalk/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
